# scripts/dance-machine.service
[Unit]
Description=Fitness Center Dance Machine Kiosk
After=network.target sound.target

[Service]
Type=simple
User=jules
WorkingDirectory=/app
ExecStart=/app/scripts/kiosk-standalone.sh
CPUSchedulingPolicy=rr
CPUSchedulingPriority=99
Nice=-20
Restart=always
RestartSec=10
Environment=DISPLAY=:0

[Install]
WantedBy=multi-user.target
