# scripts/dance-machine.service
[Unit]
Description=Fitness Center Dance Machine Kiosk
After=network.target sound.target

[Service]
Type=simple
# User=USER
# WorkingDirectory=/path/to/fcdm
Environment=DISPLAY=:0
# Systemd requires absolute paths for ExecStart when paths are used.
# Replace /path/to/fcdm with the actual deployment directory.
ExecStart=/path/to/fcdm/fcdm-orchestrator
Restart=always
RestartSec=5

# v3.1.0 Industrial Process Scheduling
CPUSchedulingPolicy=rr
CPUSchedulingPriority=99
Nice=-20
IOSchedulingClass=realtime
IOSchedulingPriority=0

[Install]
WantedBy=multi-user.target
